# deployment/supervisor-queue.conf.example
; Copy to /etc/supervisor/conf.d/mgteam-queue.conf and adjust paths/user.
; Then: sudo supervisorctl reread && sudo supervisorctl update

[program:mgteam-plataforma-queue]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/plataforma.mgteamoficial.site/artisan queue:work database --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/plataforma.mgteamoficial.site/storage/logs/queue-worker.log
stopwaitsecs=3600
